# files/lib/systemd/system/ssh.service
[Unit]
Description=ssh dummy target

[Service]
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
