# files/etc/systemd/system/unifi-core.service.d/override.conf
[Service]
MemoryMax=500M
TimeoutStartSec=300
TimeoutStopSec=300
